# 1er_año/1er_semestre/programacion/umbrello/share/autostart/kwinshutdown.desktop
[Desktop Entry]
Hidden=False
Type=Application
X-KDE-autostart-phase=0
Name=KDE session shutdown watchdog
Exec=kwinshutdown.exe --watch umbrello --timeout 5


